# telegram-alerts.service
[Unit]
Description=Telegram alerts script

[Service]
ExecStart=/opt/telegram-alerts/main.py

[Install]
WantedBy=multi-user.target
